# netscripts.timer
[Unit]
Description=Bring up network interfaces with netscripts

[Timer]
OnBootSec=1

[Install]
WantedBy=timers.target
